# .config/autostart/conky-calendar.desktop
[Desktop Entry]
Type=Application
Name=Polybar Calendar
Exec=conky -c /home/vym/.config/conky/calendar.conf
NoDisplay=true
