# .config/autostart/download.desktop
[Desktop Entry]
Type=Application
Name=Download on /TMP
Exec=mkdir /tmp/Downloads -p
OnlyShowIn=i3;
NoDisplay=true
