# .config/autostart/dunst.desktop
[Desktop Entry]
Type=Application
Name=Dunst Server
Exec=dunst -conf /home/vym/.config/dunst/dunstrc
