# .config/autostart/flameshot.desktop
[Desktop Entry]
Type=Application
Name=Flame Shot Screen
Exec=flameshot
