# .config/autostart/hacompanion.desktop
[Desktop Entry]
Type=Application
Name=Home Assistant Companion
Exec=hacompanion -config /home/vym/.config/hacompanion.toml
