# .config/autostart/polybar.desktop
[Desktop Entry]
Type=Application
Name=Polybar
Exec=polybar
OnlyShowIn=i3;
NoDisplay=true
