# .config/autostart/ulauncher.desktop
[Desktop Entry]
Type=Application
Name=Ulauncher Server
Exec=ulauncher --no-window-shadow --hide-window
NoDisplay=true
