# .config/xdg-autostart/nitrogen.desktop
[Desktop Entry]
Type=Application
Name=Nitrogen
Exec=nitrogen --set-zoom-fill --random /opt/nitrogen
OnlyShowIn=i3;
NoDisplay=true
